# backups/pipface_20260129_190337/pipface-monitor.service
[Unit]
Description=Pip Face Monitor
After=pipface.service
Requires=pipface.service

[Service]
Type=simple
WorkingDirectory=/home/nl3mos/clawd
ExecStart=/usr/bin/python3 /home/nl3mos/clawd/pip_face_monitor.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
