# config/pip_email_commands.cron
# Pip Face - Email Command Processor
# Read inbox every 5 minutes and execute commands from Nilson
*/5 * * * * /home/nl3mos/clawd/scripts/email_commands.py 10 >> /home/nl3mos/clawd/logs/email_commands.log 2>&1
